# deploy/systemd/redner-worker.service
[Unit]
Description=redner deployment worker
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/redner
EnvironmentFile=/home/ubuntu/redner/.env
ExecStart=/usr/bin/npm run start --workspace @redner/worker
Restart=on-failure
RestartSec=5
TimeoutStopSec=90

[Install]
WantedBy=multi-user.target
